# src/fty-security-wallet.service.in
[Unit]
Description=fty-security-wallet service
After=network.target malamute.service
Requires=network.target malamute.service
PartOf=bios.target

[Service]
Type=simple
User=secw-daemon
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-security-wallet -c @sysconfdir@/@PACKAGE@/fty-security-wallet.cfg

# This service is only active, as far as consumers are concerned,
# after it listens on the Unix socket used to interact with it.
# TODO: Set socket path via autoconf and template vars
### src/fty-security-wallet.cfg.in:    socket = /run/fty-security-wallet/secw.socket #   Malamute endpoint
ExecStartPost=/bin/dash -c "while sleep 0.1; do test -S /run/fty-security-wallet/secw.socket || continue; PIDS=\"`/bin/fuser /run/fty-security-wallet/secw.socket 2>/dev/null | sed -e 's,^ *,,' -e 's, *$,,'`\" && [ -n \"$PIDS\" ] || continue ; ( for P in $PIDS ; do ls -la /proc/$P/exe ; done | egrep ' @prefix@/bin/fty-security-wallet$' ) && break ; rePIDS=\"(`echo "$PIDS" | sed 's,  *,|,g'`)\" ; ps -ef | grep -v grep | grep -E \"$rePIDS\" && break ; done; echo 'Socket is being listened'"

# Last agent out should shut the door cleanly
ExecStopPost=/bin/dash -c "if [ -z \"`/bin/fuser /run/fty-security-wallet/secw.socket 2>/dev/null | sed -e 's, *,,g'`\" ] ; then rm -f /run/fty-security-wallet/secw.socket ; fi"

Restart=always

[Install]
WantedBy=bios.target
