# src/fty-security-wallet.service.in
[Unit]
Description=fty-security-wallet service
After=network.target malamute.service
Requires=network.target malamute.service
PartOf=bios.target

[Service]
Type=@systemd_service_type@
User=secw-daemon
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
EnvironmentFile=@sysconfdir@/default/bios-db-rw
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-security-wallet -c @sysconfdir@/@PACKAGE@/fty-security-wallet.cfg

# This service is only active, as far as consumers are concerned,
# after it listens on the Unix socket used to interact with it.
ExecStartPost=/bin/dash -c "while sleep 0.1; do test -S \"@socketSecurityWallet@\" || continue; PIDS=\"`/bin/fuser \"@socketSecurityWallet@\" 2>/dev/null | sed -e 's,^ *,,' -e 's, *$,,'`\" && [ -n \"$PIDS\" ] || continue ; ( for P in $PIDS ; do ls -la /proc/$P/exe ; done | egrep ' @prefix@/bin/fty-security-wallet$' ) && break ; rePIDS=\"(`echo \"$PIDS\" | sed 's,  *,|,g'`)\" ; ps -ef | grep -v grep | grep -E \"$rePIDS\" && break ; done; echo 'Socket is being listened'"

# There is also some time for malamute connection that I could not detect easily
# Perhaps we can query with an MQ client (bmsg?) or even use Type=notify unit
ExecStartPost=/bin/sleep 3

# Last agent out should shut the door cleanly
ExecStopPost=/bin/dash -c "if [ -z \"`/bin/fuser \"@socketSecurityWallet@\" 2>/dev/null | sed -e 's, *,,g'`\" ] ; then rm -f \"@socketSecurityWallet@\" ; fi"

Restart=always

[Install]
WantedBy=bios.target
