# src/fty-security-wallet.service.in
[Unit]
Description=fty-security-wallet service
After=network.target malamute.service
Requires=network.target malamute.service
PartOf=bios.target

[Service]
Type=simple
User=secw-daemon
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
EnvironmentFile=@sysconfdir@/default/bios-db-rw
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-security-wallet -c @sysconfdir@/@PACKAGE@/fty-security-wallet.cfg
Restart=always

[Install]
WantedBy=bios.target
